# supervisord.conf
[unix_http_server]
file=/tmp/supervisor_zzdh_replica_demo.sock

[supervisord]
pidfile=/tmp/supervisord_zzdh_replica_demo.pid
nodaemon=false
logfile=backend/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=5

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor_zzdh_replica_demo.sock

[program:api]
directory=/mnt/d/项目文件夹/zzdh-replica-demo/backend
command=/bin/bash -lc "source /mnt/d/项目文件夹/zzdh-replica-demo/load-env.sh && source /root/miniconda3/etc/profile.d/conda.sh && conda activate video && exec python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8000"
# 也可直接像 NexusAI 那样写在 supervisord.conf 中：
# environment=HTTP_PROXY="http://127.0.0.1:7890",HTTPS_PROXY="http://127.0.0.1:7890",ALL_PROXY="socks5://127.0.0.1:7891",NO_PROXY="localhost,127.0.0.1"
autostart=true
autorestart=true
startretries=3
priority=20
stopasgroup=true
killasgroup=true
stdout_logfile=backend/logs/supervisor_api.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=backend/logs/supervisor_api_err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
